# seed-one/wylloh.service
[Unit]
Description=Wylloh Media Platform
After=network.target

[Service]
Type=simple
User=${USER}
WorkingDirectory=/opt/wylloh
ExecStart=/usr/bin/node /opt/wylloh/index.js
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=wylloh

[Install]
WantedBy=multi-user.target 